# crontab
# Run the script every 15 minutes
*/15 * * * * root /usr/local/bin/python3 /var/youtube/script.py > /proc/1/fd/1 2>/proc/1/fd/2
